# prefect-v1-legacy/devops/infrastructure-as-code/aws/prefect-agent-on-ec2/prefect-agent.service
[Unit]
Description=Prefect Docker Agent
After=network.target
StartLimitIntervalSec=0
[Service]
Type=simple
Restart=always
RestartSec=1
User=root
ExecStart=/usr/local/bin/prefect agent docker start -k $PREFECT_API_KEY
[Install]
WantedBy=multi-user.target 